# roles/matrix-riot-web/templates/systemd/matrix-riot-web.service.j2
#jinja2: lstrip_blocks: "True"
[Unit]
Description=Matrix riot-web server
{% for service in matrix_riot_web_systemd_required_services_list %}
Requires={{ service }}
After={{ service }}
{% endfor %}

[Service]
Type=simple
ExecStartPre=/bin/sleep 5
ExecStartPre=-/usr/bin/docker-compose -f /var/lib/docker/lab/apps/matrix/docker-compose.riot.yml kill
ExecStart=/usr/bin/docker-compose -f /var/lib/docker/lab/apps/matrix/docker-compose.riot.yml up
ExecStop=-/usr/bin/docker-compose -f /var/lib/docker/lab/apps/matrix/docker-compose.riot.yml kill

Restart=always
RestartSec=30
SyslogIdentifier=matrix-riot-web

[Install]
WantedBy=multi-user.target
